# deploy/wildlens-backup.timer
[Unit]
Description=wildlens daily data backup (02:30)

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
